# etc/supervisor.d/emulator.ini
[program:emulator]
command=python3 -m http.server --directory /home/zeal8bit/Zeal-WebEmulator/ 1145
user=root
autostart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true